# centos_deploy/ana01/run_main.cron
*/5 * * * * bash /root/sharefile/ana01/run_main.sh   # 5min
# 0 * * * * bash /root/sharefile/ana01/run_main.sh   # 1hour
# 0 0 * * * bash /root/sharefile/ana01/run_main.sh   # 1day
# 0 0 * * 0 bash /root/sharefile/ana01/run_main.sh   # 1week
# 0 0 1 * * bash /root/sharefile/ana01/run_main.sh   # 1month
# 0 0 1 1 * bash /root/sharefile/ana01/run_main.sh   # 1year

